# allora-submission.service
[Unit]
Description=Allora BTC/USD Prediction Submission Daemon
Documentation=https://github.com/MISHHF93/allora-forge-builder-kit
After=network-online.target
Wants=network-online.target
Requires=network-online.target

[Service]
Type=simple
User=codespa
Group=codespa
WorkingDirectory=/workspaces/allora-forge-builder-kit
Environment="PATH=/workspaces/allora-forge-builder-kit/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Main execution
ExecStart=/workspaces/allora-forge-builder-kit/.venv/bin/python -u submit_prediction.py --daemon

# Restart policy for reliability
Restart=always
RestartSec=30
StartLimitInterval=300
StartLimitBurst=5

# Graceful shutdown: allow 30 seconds to shut down cleanly
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits (adjust as needed)
MemoryLimit=2G
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=allora-submission

# Security
PrivateTmp=yes
NoNewPrivileges=yes
ReadOnlyPaths=/etc
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
RestrictNamespaces=yes
SystemCallFilter=~@clock @debug @module @mount @obsolete @privileged @raw-io @reboot @swap
LockPersonality=yes

# Ensure clean environment
ClearEnvironment=no

# Optional: Automatically restart on system boot
WantedBy=multi-user.target

[Install]
WantedBy=multi-user.target
